# selfimprove/launchd/com.maylac.selfimprove-rot.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.maylac.selfimprove-rot</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>__HOME__/workspace/agents-setup/selfimprove/scripts/rot_check_weekly.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>__HOME__/.claude/logs/selfimprove-rot.out.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/.claude/logs/selfimprove-rot.err.log</string>
</dict>
</plist>
