# launchd/com.maylac.worktree-retention-report.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.maylac.worktree-retention-report</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>$HOME/workspace/agents-setup/scripts/worktree-retention-report.py</string>
    <string>--output</string>
    <string>$HOME/tasks/worktree-retention-report.md</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>0</integer>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>15</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>$HOME/.claude/logs/worktree-retention-report.out.log</string>
  <key>StandardErrorPath</key>
  <string>$HOME/.claude/logs/worktree-retention-report.err.log</string>
</dict>
</plist>
